# ansible/roles/kafka-connect/templates/kafka-connect.service.j2
[Unit]
Description=confluent platform kafka connect
After=network.target

[Service]
ExecStart=/usr/bin/kafka_connect_start
ExecStop=/usr/bin/kafka_connect_stop

[Install]
WantedBy=multi-user.target
